# crontab
15 2 * * 0 certbot renew --post-hook "nginx -s reload"

# Empty line left on purpose
