# crontab
15 2 * * 0 certbot renew

# Empty line left on purpose
