# crontab
15 2 * * 0 /renew-certs.sh

# Empty line left on purpose
